# systemd/fancy-frame-setup-portal.service
[Unit]
Description=Fancy Frame Wi-Fi Setup Portal
After=fancy-frame-setup-mode.service
Requires=fancy-frame-setup-mode.service

[Service]
Type=simple
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /opt/fancy-frame/portal/app.py
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
