# systemd/photo-frame-render.service
[Unit]
Description=Photo Frame Slideshow Renderer
After=local-fs.target

[Service]
User=root
Group=root
ExecStart=/opt/photo-frame/scripts/render-slideshow.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
