# systemd/photo-frame-wifi-bootstrap.service
[Unit]
Description=Photo Frame Wi-Fi Bootstrap
After=network.target

[Service]
Type=oneshot
ExecStart=/opt/photo-frame/scripts/wifi_bootstrap.sh

[Install]
WantedBy=multi-user.target
